# scheduled-agents/local/co.opsapp.nightly-bug-triage-ios.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>co.opsapp.nightly-bug-triage-ios</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-l</string>
        <string>/Users/jacksonsweet/Projects/OPS/ops-software-bible/scheduled-agents/local/nightly-bug-triage-ios.sh</string>
    </array>

    <!-- 3:15am local — staggered after the web run to reduce contention. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>15</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/jacksonsweet/Library/Logs/ops-nightly/ios-launchd.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/jacksonsweet/Library/Logs/ops-nightly/ios-launchd.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/jacksonsweet/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
